# fastapi.conf
[program:fastapi]
command=/path/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
directory=/path/to/your/app
user=fastapi
autostart=true
autorestart=true
stdout_logfile=/var/log/fastapi.log
stderr_logfile=/var/log/fastapi.err.log
